# xcal-for-conky-kouki-1.0.1/etc/skel/.config/autostart/xconky.desktop
[Desktop Entry]
Exec=/usr/bin/XCONKY
Name=xconky
Type=Application
Version=1.0
